# systemd/kanshi.service
[Unit]
Description=Dynamic display configuration

[Service]
Type=simple
ExecStart=/usr/bin/kanshi

[Install]
WantedBy=sway-session.target
