# systemd/mako.service
[Unit]
Description=Notification daemon

[Service]
Type=simple
ExecStart=/usr/bin/mako

[Install]
WantedBy=sway-session.target
