# systemd/pcloud.service
[Unit]
Description=pCloud

[Service]
Type=simple
ExecStart=/usr/bin/pcloudcc -u <email> -m %h/pcloud

[Install]
WantedBy=default.target
