# systemd/sway-session.target
[Unit]
Requires=basic.target
BindsTo=graphical-session.target
Before=graphical-session.target
Requires=basic.target
