# imwheel/imwheel.service
[Unit]
Description=imwheel daemon

[Service]
Type=simple
ExecStart=/usr/bin/imwheel -d

[Install]
WantedBy=default.target
